# crontab
#crontab

0 1 * * 1    ( date ; letsencrypt certonly"  ) >> /var/log/certs.log 2>&1
