# crontab
#crontab
0 1 * * 1    ( date ; letsencrypt certonly && nginx -t && service nginx reload  ) >> /var/log/certs.log 2>&1
