# docker-ssl-haproxy/rsyslog-supervisord.conf
[program:rsyslog]
command=rsyslogd -n -c5
autostart=true
autorestart=true
redirect_stderr=true