# ssl-proxy/haproxy-supervisord.conf
[program:haproxy]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg
autostart=true
autorestart=true
stderr_logfile=/var/log/haproxy.err.log
stdout_logfile=/var/log/haproxy.out.log


